# workspace-supervisord.conf
[unix_http_server]
file=/run/supervisor.sock
chmod=0700

[supervisord]
nodaemon=true
user=root
logfile=/var/log/mosgarage/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
loglevel=info
pidfile=/run/supervisord.pid

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///run/supervisor.sock

# ── code-server ───────────────────────────────────────────────────────────────
[program:code-server]
command=/usr/bin/code-server --config /home/mosgarage/.config/code-server/config.yaml
user=mosgarage
environment=HOME="/home/mosgarage",USER="mosgarage",SHELL="/bin/bash",DOTNET_ROOT="/usr/local/dotnet",PATH="/usr/local/dotnet:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
directory=/home/mosgarage/workspace
autostart=true
autorestart=true
startretries=5
startsecs=3
stdout_logfile=/var/log/mosgarage/code-server.log
stdout_logfile_maxbytes=5MB
stderr_logfile=/var/log/mosgarage/code-server-err.log
stderr_logfile_maxbytes=5MB
priority=10

# ── sshd ──────────────────────────────────────────────────────────────────────
[program:sshd]
command=/usr/sbin/sshd -D -e -p 2222
autostart=true
autorestart=true
startretries=5
startsecs=2
stdout_logfile=/var/log/mosgarage/sshd.log
stdout_logfile_maxbytes=5MB
stderr_logfile=/var/log/mosgarage/sshd-err.log
stderr_logfile_maxbytes=5MB
priority=5

# ── mosgarage agent ───────────────────────────────────────────────────────────
# Started only when the agent binary is present (injected by control plane).
# In standalone / WSL mode this is a no-op.
[program:mosgarage-agent]
command=/bin/bash -c 'test -f /app/agent/Mosgarage.Agent && exec /usr/local/dotnet/dotnet /app/agent/Mosgarage.Agent || sleep infinity'
user=mosgarage
environment=HOME="/home/mosgarage",DOTNET_ROOT="/usr/local/dotnet",PATH="/usr/local/dotnet:/usr/local/bin:/usr/bin:/bin",ASPNETCORE_URLS="http://+:7072"
autostart=true
autorestart=true
startretries=3
startsecs=5
stdout_logfile=/var/log/mosgarage/agent.log
stdout_logfile_maxbytes=5MB
stderr_logfile=/var/log/mosgarage/agent-err.log
stderr_logfile_maxbytes=5MB
priority=20
