# config/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/mosgarage/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info

; ── git-setup (runs once at boot — clones/pulls repo) ────────
[program:git-setup]
command=/usr/local/bin/git-setup
user=mosgarage
environment=HOME="/home/mosgarage",USER="mosgarage",GITHUB_USER="%(ENV_GITHUB_USER)s",GITHUB_TOKEN="%(ENV_GITHUB_TOKEN)s",GITHUB_ORG="%(ENV_GITHUB_ORG)s",GITHUB_REPO="%(ENV_GITHUB_REPO)s",GIT_BRANCH="%(ENV_GIT_BRANCH)s",GIT_NAME="%(ENV_GIT_NAME)s",GIT_EMAIL="%(ENV_GIT_EMAIL)s"
autostart=true
autorestart=false
startsecs=0
startretries=3
priority=10
stdout_logfile=/var/log/mosgarage/git-sync.log
stderr_logfile=/var/log/mosgarage/git-sync.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=0

; ── git-sync (continuous backup daemon) ──────────────────────
[program:git-sync]
command=/usr/local/bin/git-sync
user=mosgarage
environment=HOME="/home/mosgarage",USER="mosgarage",GITHUB_USER="%(ENV_GITHUB_USER)s",GITHUB_TOKEN="%(ENV_GITHUB_TOKEN)s",GITHUB_ORG="%(ENV_GITHUB_ORG)s",GITHUB_REPO="%(ENV_GITHUB_REPO)s",GIT_BRANCH="%(ENV_GIT_BRANCH)s",GIT_NAME="%(ENV_GIT_NAME)s",GIT_EMAIL="%(ENV_GIT_EMAIL)s",GIT_SYNC_INTERVAL="%(ENV_GIT_SYNC_INTERVAL)s"
autostart=true
autorestart=true
startretries=10
startsecs=5
priority=20
stdout_logfile=/var/log/mosgarage/git-sync.log
stderr_logfile=/var/log/mosgarage/git-sync.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=0

; ── code-server (VS Code in browser) ─────────────────────────
[program:code-server]
command=code-server --config /home/mosgarage/.config/code-server/config.yaml /app/workspace
user=mosgarage
environment=HOME="/home/mosgarage",USER="mosgarage",PASSWORD="%(ENV_CODE_SERVER_PASSWORD)s"
autostart=true
autorestart=true
startretries=5
startsecs=3
priority=30
stdout_logfile=/var/log/mosgarage/code-server.log
stderr_logfile=/var/log/mosgarage/code-server.err
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB

; ── Node home server ──────────────────────────────────────────
[program:node-server]
command=node /app/server/index.js
user=mosgarage
environment=HOME="/home/mosgarage",USER="mosgarage",NODE_ENV="%(ENV_NODE_ENV)s",PORT="%(ENV_NODE_SERVER_PORT)s"
directory=/app/server
autostart=true
autorestart=true
startretries=5
startsecs=2
priority=40
stdout_logfile=/var/log/mosgarage/node-server.log
stderr_logfile=/var/log/mosgarage/node-server.err
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB

; ── Express API ───────────────────────────────────────────────
[program:api-server]
command=node /app/api/index.js
user=mosgarage
environment=HOME="/home/mosgarage",USER="mosgarage",NODE_ENV="%(ENV_NODE_ENV)s",PORT="%(ENV_API_PORT)s"
directory=/app/api
autostart=true
autorestart=true
startretries=5
startsecs=2
priority=40
stdout_logfile=/var/log/mosgarage/api-server.log
stderr_logfile=/var/log/mosgarage/api-server.err
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB

; ── Supervisor internal ───────────────────────────────────────
[inet_http_server]
port=127.0.0.1:9001

[supervisorctl]
serverurl=http://127.0.0.1:9001

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface
